# deploy/jt-glogarch.service
[Unit]
Description=jt-glogarch - Graylog Open Archive
Documentation=https://github.com/jasoncheng7115/jt-glogarch
After=network.target

[Service]
Type=simple
User=jt-glogarch
Group=jt-glogarch
WorkingDirectory=/opt/jt-glogarch
ExecStart=/usr/local/bin/glogarch server
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

# Memory safety cap (jt-glogarch is usually co-located with Graylog + OpenSearch on
# ONE VM). The archive reader streams and the import has a memory guard, so normal
# usage stays well under 1.5 GB. This is a SOFT cap: past MemoryHigh the kernel
# throttles jt-glogarch and reclaims its memory (protecting OpenSearch/Graylog and
# the box) but NEVER OOM-kills it — deliberately no hard limit, so heavy PDF
# reports (Chromium runs in this same cgroup) are never killed. Raise on large boxes.
MemoryAccounting=yes
MemoryHigh=3G

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/data/graylog-archives /opt/jt-glogarch

[Install]
WantedBy=multi-user.target
