# deploy/jt-glogarch.service
[Unit]
Description=jt-glogarch - Graylog Open Archive
Documentation=https://github.com/jasoncheng7115/jt-glogarch
After=network.target

[Service]
Type=simple
User=jt-glogarch
Group=jt-glogarch
WorkingDirectory=/opt/jt-glogarch
ExecStart=/usr/local/bin/glogarch server
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

# Memory safety cap (jt-glogarch is usually co-located with Graylog + OpenSearch on
# ONE VM). The archive reader streams and the import has a memory guard, so normal
# usage stays well under 1.5 GB — but this cgroup cap is a hard backstop: if
# jt-glogarch ever ran away it is throttled at MemoryHigh and, at MemoryMax, the
# kernel OOM-kills ONLY this service's cgroup — never OpenSearch/Graylog or the box.
# Raise these on large boxes or heavy PDF-report use; lower them on very small VMs.
MemoryAccounting=yes
MemoryHigh=3G
MemoryMax=4G

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/data/graylog-archives /opt/jt-glogarch

[Install]
WantedBy=multi-user.target
